# systemd/telegram-codex-controller.user.service
[Unit]
Description=Private Telegram tmux/Codex controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/telegram-tmux-control.env
Environment=TELEGRAM_TMUX_CONFIG=%h/.config/telegram-tmux-control.env
Environment=TELEGRAM_TMUX_STATE=%h/.local/state/telegram-tmux-control.json
Environment=TELEGRAM_TMUX_SESSION=web
Environment=TELEGRAM_TMUX_WORKSPACE=%h
ExecStart=%h/.local/share/telegram-codex-controller/telegram-tmux-control.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
