# systemd/telegram-tmux-controller.service.example
[Unit]
Description=Private Telegram tmux/Codex controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace every YOUR_USER and /path/to value before installing as a system service.
# For the simpler per-user service, use ../install.sh instead.
User=YOUR_USER
WorkingDirectory=/path/to/telegram-codex-controller
Environment=TELEGRAM_TMUX_CONFIG=/home/YOUR_USER/.config/telegram-tmux-control.env
Environment=TELEGRAM_TMUX_STATE=/home/YOUR_USER/.local/state/telegram-tmux-control.json
Environment=TELEGRAM_TMUX_SESSION=web
Environment=TELEGRAM_TMUX_WORKSPACE=/path/to/workspace
Environment=TELEGRAM_TMUX_CODEX_BIN=/path/to/codex
Environment=TELEGRAM_TMUX_NODE_BIN=/path/to/node
ExecStart=/path/to/telegram-codex-controller/telegram-tmux-control.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
